# antix-livefs/var/lib/dpkg/info/bootlogd.list
/.
/etc
/etc/init.d
/etc/init.d/bootlogd
/etc/init.d/stop-bootlogd
/etc/init.d/stop-bootlogd-single
/usr
/usr/bin
/usr/bin/readbootlog
/usr/sbin
/usr/sbin/bootlogd
/usr/share
/usr/share/doc
/usr/share/doc/bootlogd
/usr/share/doc/bootlogd/changelog.Debian.gz
/usr/share/doc/bootlogd/changelog.gz
/usr/share/doc/bootlogd/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/bootlogd
/usr/share/man
/usr/share/man/de
/usr/share/man/de/man1
/usr/share/man/de/man1/readbootlog.1.gz
/usr/share/man/de/man8
/usr/share/man/de/man8/bootlogd.8.gz
/usr/share/man/es
/usr/share/man/es/man1
/usr/share/man/es/man1/readbootlog.1.gz
/usr/share/man/es/man8
/usr/share/man/es/man8/bootlogd.8.gz
/usr/share/man/fi
/usr/share/man/fi/man1
/usr/share/man/fi/man1/readbootlog.1.gz
/usr/share/man/fi/man8
/usr/share/man/fi/man8/bootlogd.8.gz
/usr/share/man/fr
/usr/share/man/fr/man1
/usr/share/man/fr/man1/readbootlog.1.gz
/usr/share/man/fr/man8
/usr/share/man/fr/man8/bootlogd.8.gz
/usr/share/man/hu
/usr/share/man/hu/man1
/usr/share/man/hu/man1/readbootlog.1.gz
/usr/share/man/hu/man8
/usr/share/man/hu/man8/bootlogd.8.gz
/usr/share/man/id
/usr/share/man/id/man1
/usr/share/man/id/man1/readbootlog.1.gz
/usr/share/man/id/man8
/usr/share/man/id/man8/bootlogd.8.gz
/usr/share/man/man1
/usr/share/man/man1/readbootlog.1.gz
/usr/share/man/man8
/usr/share/man/man8/bootlogd.8.gz
/usr/share/man/pl
/usr/share/man/pl/man1
/usr/share/man/pl/man1/readbootlog.1.gz
/usr/share/man/pl/man8
/usr/share/man/pl/man8/bootlogd.8.gz
/usr/share/man/pt
/usr/share/man/pt/man1
/usr/share/man/pt/man1/readbootlog.1.gz
/usr/share/man/pt/man8
/usr/share/man/pt/man8/bootlogd.8.gz
