# antix-livefs/var/lib/dpkg/info/check-dfsg-status.list
/.
/etc
/etc/cron.monthly
/etc/cron.monthly/check-dfsg-status
/usr
/usr/bin
/usr/bin/check-dfsg-status
/usr/share
/usr/share/check-dfsg-status
/usr/share/check-dfsg-status/reasons
/usr/share/check-dfsg-status/reasons/check-dfsg-status
/usr/share/doc
/usr/share/doc/check-dfsg-status
/usr/share/doc/check-dfsg-status/README.Debian
/usr/share/doc/check-dfsg-status/changelog.gz
/usr/share/doc/check-dfsg-status/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/check-dfsg-status.1.gz
