# antix-livefs/var/lib/dpkg/info/desktop-defaults-core-antix.list
/.
/etc
/etc/X11
/etc/X11/xorg.conf.d
/etc/X11/xorg.conf.d/synaptics.conf
/etc/X11/xorg.conf.in
/etc/init.d
/etc/init.d/umountnfs-alternative.sh
/etc/modprobe.d
/etc/modprobe.d/hang-on-shutdown.conf
/etc/modprobe.d/i915-power-saving.conf
/etc/modprobe.d/iwlwifi-no-blink.conf
/etc/skel
/etc/skel/.Xresources
/etc/skel/.config
/etc/skel/.config/htop
/etc/skel/.config/htop/htoprc
/etc/skel/.config/xresources.d
/etc/skel/.config/xresources.d/xcalc
/etc/skel/.config/xresources.d/xcalc/xcalc.conf
/etc/skel/.config/xresources.d/xcolours
/etc/skel/.config/xresources.d/xcolours/antix.conf
/etc/skel/.config/xresources.d/xcolours/zenburn.conf
/etc/skel/.conkyrc
/etc/skel/.conkyrc-lua
/etc/skel/.nanorc
/etc/skel/.xinitrc-custom
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/61-optical-polling-rules.rules
/etc/udev/rules.d/80-net-name-slot.rules
/etc/udev/rules.d/80-net-setup-link.rules
/etc/udev/rules.d/90-fstab-automount.rules
/etc/udev/rules.d/90-usb-semiauto.rules
/root
/root/.config
/root/.config/htop
/root/.config/htop/htoprc
/root/.nanorc
/usr
/usr/local
/usr/local/bin
/usr/local/bin/conky-colors
/usr/sbin
/usr/sbin/autologin
/usr/share
/usr/share/X11
/usr/share/X11/xkb
/usr/share/X11/xkb/symbols
/usr/share/X11/xkb/symbols/br-abnt2
/usr/share/X11/xorg.conf.d
/usr/share/X11/xorg.conf.d/20-thinkpad.conf
/usr/share/doc
/usr/share/doc/desktop-defaults-core-antix
/usr/share/doc/desktop-defaults-core-antix/changelog.gz
/usr/share/doc/desktop-defaults-core-antix/copyright
/usr/share/keymaps
/usr/share/keymaps/i386
/usr/share/keymaps/i386/azerty
/usr/share/keymaps/i386/qwerty
/usr/share/keymaps/i386/qwertz
/usr/share/keymaps/i386/azerty/be.kmap.gz
/usr/share/keymaps/i386/qwerty/br.kmap.gz
/usr/share/keymaps/i386/qwerty/cz.kmap.gz
/usr/share/keymaps/i386/qwerty/gb.kmap.gz
/usr/share/keymaps/i386/qwerty/is.kmap.gz
/usr/share/keymaps/i386/qwerty/jp.kmap.gz
/usr/share/keymaps/i386/qwerty/lv.kmap.gz
/usr/share/keymaps/i386/qwerty/pt.kmap.gz
/usr/share/keymaps/i386/qwerty/se.kmap.gz
/usr/share/keymaps/i386/qwerty/sk.kmap.gz
/usr/share/keymaps/i386/qwerty/tr.kmap.gz
/usr/share/keymaps/i386/qwertz/hr.kmap.gz
/usr/share/keymaps/i386/qwertz/si.kmap.gz
