# antix-livefs/var/lib/dpkg/info/init-diversity-tools.list
/.
/etc
/etc/init.d
/usr
/usr/lib
/usr/lib/init-diversity
/usr/lib/init-diversity/poweroff.sh
/usr/lib/init-diversity/rc.sh
/usr/lib/init-diversity/rcS.sh
/usr/lib/init-diversity/reboot.sh
/usr/lib/init-diversity/runlevel.sh
/usr/lib/init-diversity/shutdown.sh
/usr/lib/init-diversity/user-session.sh
/usr/sbin
/usr/share
/usr/share/doc
/usr/share/doc/init-diversity-tools
/usr/share/doc/init-diversity-tools/changelog.gz
/usr/share/doc/init-diversity-tools/copyright
/usr/share/init-diversity-tools
/usr/share/init-diversity-tools/diversions
/usr/share/init-diversity-tools/diversions/etc
/usr/share/init-diversity-tools/diversions/etc/init.d
/usr/share/init-diversity-tools/diversions/usr
/usr/share/init-diversity-tools/diversions/usr/sbin
/usr/share/init-diversity-tools/poweroff.sh
/usr/share/init-diversity-tools/rc.sh
/usr/share/init-diversity-tools/rcS.sh
/usr/share/init-diversity-tools/reboot.sh
/usr/share/init-diversity-tools/runlevel.sh
/usr/share/init-diversity-tools/shutdown.sh
/usr/share/init-diversity-tools/user-session.sh
/etc/init.d/rc
/etc/init.d/rcS
/usr/sbin/poweroff
/usr/sbin/reboot
/usr/sbin/runlevel
/usr/sbin/shutdown
