# antix-livefs/var/lib/dpkg/info/plocate.conffiles
/etc/cron.daily/plocate
/etc/updatedb.conf
