# antix-livefs/usr/share/s6-rc/services/keyboard-setup/up
#!/usr/bin/execlineb -P

exec sh /lib/console-setup/keyboard-setup.sh
