# antix-livefs/usr/share/s6-rc/services/setup-loopback/up
#!/usr/bin/execlineb -P

exec -c ip link set up dev lo
