# antix-livefs/etc/s6-linux-init/current-release/run-image/service/s6-svscan-log/run
#!/usr/bin/execlineb -P

/usr/lib/execline/bin/fdmove -c 1 2
/usr/lib/execline/bin/redirfd -rnb 0 fifo
/usr/bin/s6-log -bpd3 -- 1 t /run/uncaught-logs
