# antix-livefs/etc/sv/getty-tty4/finish
#!/bin/sh -euC

if [ "$1" = 163 ] ; then
    echo "runsv: ERROR in ${PWD}: another getty already on tty4"
    exec sv d .
fi
exec utmpset -w tty4
