# antix-livefs/etc/s6-rc/current/run-image/service/s6-linux-init-shutdownd/run
#!/usr/bin/execlineb -P

/usr/bin/s6-linux-init-shutdownd -d3 -c "/etc/s6-rc/current" -g 3000
