# antix-livefs/usr/share/s6-rc/services/mount-rw/up
#!/usr/bin/execlineb -P
foreground { mount -o remount,rw / }
foreground { sh /etc/s6-rc/scripts/00-shared/log-service-reload.sh }
