# antix-livefs/usr/share/s6-rc/services/mount-netdev/up
#!/usr/bin/execlineb -P

exec -c mount -a -O _netdev
