# antix-livefs/etc/cron.monthly/check-dfsg-status
#!/bin/sh

[ -x /usr/bin/check-dfsg-status ] || exit 0

/usr/bin/check-dfsg-status --quiet
