# antix-livefs/etc/sv/getty-ttyv0/finish
#!/bin/sh -euC
SERVICE=getty-ttyv0
if [ "$1" = 163 ] ; then
    echo "runsv: ERROR in $SERVICE: another getty already on ttyv0"
    sv d "$(dirname "${0}")"
    exit 0
fi
exec utmpset -w ttyv0
