# antix-livefs/usr/share/s6-rc/services/console-setup/up
#!/usr/bin/execlineb -P

exec sh /lib/console-setup/console-setup.sh
