# antix-livefs/usr/share/s6-rc/services/mount-swap/up
#!/usr/bin/execlineb -P

exec -c swapon -a
