# antix-livefs/etc/udev/rules.d/99-ipv6-privacy.rules
ACTION=="add", SUBSYSTEM=="net", RUN+="/usr/sbin/sysctl -w net.ipv6.conf.%k.use_tempaddr=2"
