# antix-livefs/usr/share/s6-rc/services/udevadm/up
#!/usr/bin/execlineb -P

foreground { redirfd -w 1 /dev/null /bin/udevadm trigger --action=add --type=subsystems }
foreground { redirfd -w 1 /dev/null /bin/udevadm trigger --action=add --type=devices }
foreground { redirfd -w 1 /dev/null /bin/udevadm settle }
