# antix-livefs/usr/share/s6-rc/services/seatd-srv/run
#!/usr/bin/execlineb -P

envfile /etc/s6-rc/config/seatd.conf
importas -sCuD "" Args Args

fdmove -c 2 1
exec -c seatd ${Args}
