# antix-livefs/usr/share/s6-rc/services/modules/up
#!/usr/bin/execlineb -P

importas -sCuD "/usr/sbin:/usr/bin:/sbin:/bin" PATH PATH

exec sh /etc/s6-rc/scripts/modules/modules.sh
