# antix-livefs/etc/s6-rc/current/run-image/service/s6-linux-init-early-getty/run
#!/usr/bin/execlineb -P

/sbin/agetty tty12 38400 linux
